# com.john.telegrambot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.john.telegrambot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Joel/REPOS/BookieGrabber/venv/bin/python</string>
        <string>/Users/Joel/REPOS/BookieGrabber/telegram_alerts.py</string>
        <string>bot</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/Joel/REPOS/BookieGrabber</string>

    <key>StandardOutPath</key>
    <string>/Users/Joel/REPOS/BookieGrabber/logs/telegrambot.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/Joel/REPOS/BookieGrabber/logs/telegrambot.err</string>

</dict>
</plist>