# com.john.bookiepipeline.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Label for your job -->
    <key>Label</key>
    <string>com.john.bookiepipeline.hourly</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/Joel/REPOS/BookieGrabber/venv/bin/python</string>
        <string>/Users/Joel/REPOS/BookieGrabber/bookie_grabber.py</string>
    </array>

    <!-- Run every 30 minutes -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Run immediately on load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Do NOT try to restart if it crashes -->
    <key>KeepAlive</key>
    <false/>

    <!-- Set working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/Joel/REPOS/BookieGrabber</string>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/Users/Joel/REPOS/BookieGrabber/logs/bookiepipeline.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/Joel/REPOS/BookieGrabber/logs/bookiepipeline.err</string>

</dict>
</plist>
