# com.john.optimizer.monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.john.optimizer.monthly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Joel/REPOS/BookieGrabber/venv/bin/python</string>
        <string>/Users/Joel/REPOS/BookieGrabber/optimizer/run_and_email.py</string>
    </array>

    <!-- Run on the 15th of each month at 08:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>15</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/Joel/REPOS/BookieGrabber</string>

    <key>StandardOutPath</key>
    <string>/Users/Joel/REPOS/BookieGrabber/logs/optimizer.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/Joel/REPOS/BookieGrabber/logs/optimizer.err</string>

</dict>
</plist>